# examples/com.openclaw.siri-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.siri-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/path/to/openclaw-voice-commands/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>API_KEY</key>
        <string>your-api-key-here</string>
        <key>CLI_COMMAND</key>
        <string>clawdbot</string>
        <key>BIND_HOST</key>
        <string>127.0.0.1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/siri-bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/siri-bridge.log</string>
</dict>
</plist>
